# supervisord.conf
[supervisord]
nodaemon=true

[program:postfix]
command=sh /opt/postfix.sh

[program:syslog-ng]
command=/usr/sbin/syslog-ng --foreground --no-caps
stdout_syslog=true

